# s6/gcsproxy/run
#!/usr/bin/env bash

set -e

exec gcsproxy -v
