# s6/nginx/run
#!/usr/bin/env bash

set -e

source /panubo-functions.sh

wait_tcp metadata.google.internal 80 5 5

TEMPLATES=(
  "/etc/nginx/conf.d/default.conf.tmpl"
  "/etc/nginx/common-locations.conf.tmpl"
)

curl -fSs -H "Metadata-Flavor: Google" http://metadata.google.internal/computeMetadata/v1/instance/attributes/pnb-proxy -o /pnb-proxy.json

for item in "${TEMPLATES[@]}"; do
  gomplate -d pnb_proxy=/pnb-proxy.json < "${item}" > "${item/%\.tmpl/}"
  echo "==> ${item/%\.tmpl/} <=="
  cat "${item/%\.tmpl/}"
done

exec nginx -g "daemon off;"
